# nvme-pwr-ctrl/xyz.openbmc_project.Control.Nvme.Power.service.in
[Unit]
Description=NVMe Power Control

[Service]
Restart=always
RestartSec=3
StartLimitBurst=10
ExecStart=@bindir@/nvme_pwr_ctrl
SyslogIdentifier=nvme_pwr_ctrl
Type=dbus
BusName=xyz.openbmc_project.Control.Nvme.Power

[Install]
WantedBy=multi-user.target
